# uwsgi-init.conf
# /etc/init/uwsgi.conf
start on startup
task
exec uwsgi --ini /etc/uwsgi/uwsgi-emperor.ini
